# var/lib/opkg/info/alsa-state.prerm
#!/bin/sh
# Begin section update-rc.d
if [ -z "$D" -a -x "/etc/init.d/alsa-state" ]; then
	/etc/init.d/alsa-state stop || :
fi
# End section update-rc.d
