# var/lib/opkg/info/arago-gpl-notice.list
/lib/systemd/system/gplv3-notice.service
/etc/init.d/gplv3-notice
