# var/lib/opkg/info/arago-gpl-notice.prerm
#!/bin/sh
OPTS=""

if [ -n "$D" ]; then
    OPTS="--root=$D"
fi

if type systemctl >/dev/null 2>/dev/null; then
	if [ -z "$D" ]; then
		systemctl stop gplv3-notice.service
	fi

	systemctl $OPTS disable gplv3-notice.service
fi
# Begin section update-rc.d
if [ -z "$D" -a -x "/etc/init.d/gplv3-notice" ]; then
	/etc/init.d/gplv3-notice stop || :
fi
# End section update-rc.d
