# var/lib/opkg/info/bt-enable.list
/lib/systemd/system/bt-enable.service
/etc/init.d/bt-enable.sh
