# var/lib/opkg/info/busybox-syslog.list
/lib/systemd/
/lib/systemd/system/
/etc/syslog.conf.busybox
/etc/syslog-startup.conf.busybox
/lib/systemd/system/busybox-syslog.service
/etc/default/busybox-syslog
/lib/systemd/system/syslog.service
/lib/systemd/system/busybox-klogd.service
/etc/init.d/syslog.busybox
