# var/lib/opkg/info/gdbc6x.prerm
#!/bin/sh
# Begin section update-rc.d
if [ -z "$D" -a -x "/etc/init.d/gdbserverproxy" ]; then
	/etc/init.d/gdbserverproxy stop || :
fi
# End section update-rc.d
