# var/lib/opkg/info/gdbc6x.list
/usr/bin/gdbc6x
/etc/init.d/gdbserverproxy
