# var/lib/opkg/info/hostap-daemon-wl18xx.prerm
#!/bin/sh
# Begin section update-rc.d
if [ -z "$D" -a -x "/etc/init.d/hostapd" ]; then
	/etc/init.d/hostapd stop || :
fi
# End section update-rc.d
