# var/lib/opkg/info/initscript-telnetd.prerm
#!/bin/sh
# Begin section update-rc.d
if [ -z "$D" -a -x "/etc/init.d/telnetd" ]; then
	/etc/init.d/telnetd stop || :
fi
# End section update-rc.d
