# var/lib/opkg/info/kdump.list
/usr/sbin/kdump
/etc/init.d/kdump
/etc/sysconfig/kdump.conf
/etc/sysconfig/
