# var/lib/opkg/info/kdump.prerm
#!/bin/sh
# Begin section update-rc.d
if [ -z "$D" -a -x "/etc/init.d/kdump" ]; then
	/etc/init.d/kdump stop || :
fi
# End section update-rc.d
