# var/lib/opkg/info/modutils-initscripts.prerm
#!/bin/sh
# Begin section update-rc.d
if [ -z "$D" -a -x "/etc/init.d/modutils.sh" ]; then
	/etc/init.d/modutils.sh stop || :
fi
# End section update-rc.d
