# var/lib/opkg/info/nfs-utils-client.prerm
#!/bin/sh
# Begin section update-rc.d
if [ -z "$D" -a -x "/etc/init.d/nfscommon" ]; then
	/etc/init.d/nfscommon stop || :
fi
# End section update-rc.d
