# var/lib/opkg/info/ptpd.list
/lib/systemd/system/ptpd.service
/usr/bin/ptpd2
/etc/default/ptpd
