# var/lib/opkg/info/rng-tools.list
/etc/init.d/rng-tools
/etc/default/rng-tools
/usr/bin/rngtest
/usr/sbin/rngd
