# var/lib/opkg/info/strongswan.list
/usr/share/strongswan/templates/config/strongswan.d/charon-logging.conf
/usr/lib/ipsec/plugins/libstrongswan-pkcs1.so
/etc/strongswan.d/starter.conf
/usr/lib/ipsec/plugins/libstrongswan-dnskey.so
/usr/lib/ipsec/plugins/
/usr/share/strongswan/templates/
/usr/share/strongswan/templates/config/plugins/fips-prf.conf
/usr/lib/ipsec/plugins/libstrongswan-xauth-generic.so
/usr/share/strongswan/templates/config/plugins/kernel-netlink.conf
/etc/strongswan.d/charon/pubkey.conf
/etc/strongswan.d/pki.conf
/usr/lib/ipsec/plugins/libstrongswan-revocation.so
/usr/lib/ipsec/plugins/libstrongswan-pem.so
/usr/share/strongswan/templates/config/
/usr/share/strongswan/templates/config/plugins/attr.conf
/usr/share/strongswan/templates/config/plugins/cmac.conf
/etc/strongswan.d/charon/pkcs11.conf
/usr/lib/ipsec/plugins/libstrongswan-vici.so
/usr/share/strongswan/templates/config/plugins/openssl.conf
/usr/lib/ipsec/libstrongswan.so.0
/usr/share/strongswan/templates/config/plugins/sshkey.conf
/usr/lib/ipsec/plugins/libstrongswan-kernel-netlink.so
/usr/share/strongswan/templates/config/plugins/aes.conf
/usr/share/strongswan/templates/config/strongswan.conf
/etc/strongswan.d/charon/openssl.conf
/usr/share/strongswan/templates/config/plugins/sha1.conf
/etc/ipsec.d/aacerts/
/usr/lib/ipsec/plugins/libstrongswan-md5.so
/usr/lib/ipsec/plugins/libstrongswan-fips-prf.so
/etc/strongswan.d/charon/aes.conf
/usr/lib/ipsec/plugins/libstrongswan-pkcs11.so
/usr/lib/ipsec/libstrongswan.so.0.0.0
/usr/share/strongswan/templates/config/strongswan.d/starter.conf
/usr/share/strongswan/templates/config/plugins/x509.conf
/usr/share/strongswan/
/etc/ipsec.d/acerts/
/etc/ipsec.conf
/etc/strongswan.d/charon/updown.conf
/usr/share/strongswan/templates/config/plugins/ctr.conf
/usr/share/strongswan/templates/config/plugins/pubkey.conf
/usr/lib/ipsec/libcharon.so.0.0.0
/etc/strongswan.d/charon/pkcs12.conf
/etc/strongswan.d/charon/xauth-generic.conf
/usr/share/strongswan/templates/config/plugins/rc2.conf
/usr/lib/ipsec/plugins/libstrongswan-pkcs12.so
/usr/lib/ipsec/plugins/libstrongswan-aes.so
/usr/lib/ipsec/plugins/libstrongswan-hmac.so
/etc/strongswan.d/charon/fips-prf.conf
/etc/strongswan.d/charon/pkcs7.conf
/etc/strongswan.d/charon/ctr.conf
/usr/bin/pki
/etc/strongswan.d/charon/revocation.conf
/usr/lib/ipsec/plugins/libstrongswan-pubkey.so
/usr/lib/ipsec/libvici.so.0.0.0
/usr/share/strongswan/templates/config/strongswan.d/pki.conf
/usr/share/strongswan/templates/config/plugins/vici.conf
/usr/share/strongswan/templates/config/plugins/sha2.conf
/etc/strongswan.d/charon/rc2.conf
/usr/lib/ipsec/plugins/libstrongswan-socket-default.so
/usr/share/strongswan/templates/config/plugins/pkcs11.conf
/etc/ipsec.d/certs/
/usr/share/strongswan/templates/config/plugins/pgp.conf
/usr/lib/ipsec/plugins/libstrongswan-cmac.so
/etc/strongswan.d/charon/pkcs1.conf
/usr/lib/ipsec/plugins/libstrongswan-x509.so
/usr/lib/ipsec/plugins/libstrongswan-stroke.so
/etc/ipsec.d/reqs/
/usr/share/strongswan/templates/config/plugins/hmac.conf
/etc/strongswan.d/charon/nonce.conf
/usr/share/strongswan/templates/config/plugins/xcbc.conf
/etc/strongswan.d/charon/attr.conf
/etc/ipsec.d/
/etc/ipsec.d/crls/
/etc/ipsec.d/ocspcerts/
/etc/strongswan.d/charon/pgp.conf
/etc/strongswan.d/charon/cmac.conf
/etc/strongswan.conf
/usr/lib/ipsec/plugins/libstrongswan-ctr.so
/etc/strongswan.d/
/usr/libexec/ipsec/charon
/etc/strongswan.d/charon/pkcs8.conf
/etc/strongswan.d/charon/sha1.conf
/usr/lib/ipsec/
/usr/share/strongswan/templates/config/plugins/updown.conf
/usr/share/strongswan/templates/config/plugins/pkcs12.conf
/usr/share/strongswan/templates/config/plugins/curl.conf
/usr/share/strongswan/templates/config/plugins/constraints.conf
/usr/share/strongswan/templates/config/plugins/pem.conf
/etc/strongswan.d/charon/socket-default.conf
/etc/strongswan.d/charon/random.conf
/usr/libexec/ipsec/_copyright
/usr/share/strongswan/templates/config/plugins/revocation.conf
/usr/lib/ipsec/plugins/libstrongswan-pkcs7.so
/etc/strongswan.d/charon/x509.conf
/etc/strongswan.d/charon/constraints.conf
/etc/strongswan.d/charon/pem.conf
/usr/lib/ipsec/plugins/libstrongswan-rc2.so
/etc/strongswan.d/charon/dnskey.conf
/usr/lib/ipsec/libcharon.so.0
/etc/strongswan.d/pkcs11_plugin.conf
/usr/libexec/ipsec/
/etc/strongswan.d/charon/stroke.conf
/usr/lib/ipsec/plugins/libstrongswan-curl.so
/etc/strongswan.d/charon/vici.conf
/etc/strongswan.d/charon/sha2.conf
/etc/strongswan.d/charon/sqlite.conf
/etc/ipsec.d/cacerts/
/usr/lib/ipsec/plugins/libstrongswan-pgp.so
/usr/lib/ipsec/plugins/libstrongswan-pkcs8.so
/usr/lib/ipsec/libvici.so.0
/lib/systemd/system/strongswan.service
/etc/strongswan.d/charon/hmac.conf
/usr/libexec/ipsec/starter
/etc/strongswan.d/charon/xcbc.conf
/usr/share/strongswan/templates/config/strongswan.d/
/usr/lib/ipsec/plugins/libstrongswan-resolve.so
/etc/ipsec.d/private/
/usr/lib/ipsec/plugins/libstrongswan-updown.so
/usr/lib/ipsec/plugins/libstrongswan-random.so
/usr/lib/ipsec/plugins/libstrongswan-openssl.so
/usr/share/strongswan/templates/config/plugins/socket-default.conf
/usr/share/strongswan/templates/config/plugins/random.conf
/usr/share/strongswan/templates/config/plugins/xauth-generic.conf
/usr/libexec/ipsec/stroke
/usr/lib/ipsec/plugins/libstrongswan-sqlite.so
/usr/share/strongswan/templates/config/plugins/pkcs7.conf
/usr/lib/ipsec/plugins/libstrongswan-sshkey.so
/usr/lib/ipsec/plugins/libstrongswan-des.so
/etc/strongswan.d/charon/curl.conf
/usr/share/strongswan/templates/config/plugins/dnskey.conf
/usr/share/strongswan/templates/config/plugins/stroke.conf
/usr/share/strongswan/templates/config/plugins/des.conf
/usr/sbin/ipsec
/usr/share/strongswan/templates/config/plugins/pkcs1.conf
/etc/strongswan.d/charon/kernel-netlink.conf
/usr/share/strongswan/templates/config/plugins/sqlite.conf
/usr/share/strongswan/templates/config/plugins/nonce.conf
/etc/strongswan.d/charon.conf
/usr/share/strongswan/templates/config/plugins/resolve.conf
/usr/share/strongswan/templates/config/strongswan.d/charon.conf
/usr/lib/ipsec/plugins/libstrongswan-sha1.so
/etc/strongswan.d/charon/des.conf
/usr/share/strongswan/templates/config/plugins/md5.conf
/usr/share/strongswan/templates/config/plugins/
/usr/lib/ipsec/plugins/libstrongswan-xcbc.so
/usr/libexec/ipsec/_updown
/etc/strongswan.d/charon/resolve.conf
/etc/strongswan.d/charon/md5.conf
/usr/lib/ipsec/plugins/libstrongswan-attr.so
/etc/strongswan.d/charon/
/usr/lib/ipsec/plugins/libstrongswan-nonce.so
/usr/share/strongswan/templates/config/plugins/pkcs8.conf
/etc/strongswan.d/charon/sshkey.conf
/usr/lib/ipsec/plugins/libstrongswan-sha2.so
/etc/strongswan.d/charon-logging.conf
/usr/lib/ipsec/plugins/libstrongswan-constraints.so
