# var/lib/opkg/info/thermal-init.prerm
#!/bin/sh
# Begin section update-rc.d
if [ -z "$D" -a -x "/etc/init.d/thermal-zone-init" ]; then
	/etc/init.d/thermal-zone-init stop || :
fi
# End section update-rc.d
