# var/lib/opkg/info/thttpd.prerm
#!/bin/sh
# Begin section update-rc.d
if [ -z "$D" -a -x "/etc/init.d/thttpd" ]; then
	/etc/init.d/thttpd stop || :
fi
# End section update-rc.d
