# var/lib/opkg/info/uim.prerm
#!/bin/sh
# Begin section update-rc.d
if [ -z "$D" -a -x "/etc/init.d/uim-sysfs" ]; then
	/etc/init.d/uim-sysfs stop || :
fi
# End section update-rc.d
