# var/lib/opkg/info/weston-init.list
/usr/bin/runWeston
/etc/init.d/weston
/usr/bin/weston-start
